# openclaw-imp.service
[Unit]
Description=Imp Zero OpenClaw Voice Terminal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/openclaw-imp
EnvironmentFile=/home/pi/openclaw-imp/.env
ExecStart=/usr/bin/python3 main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
